# ssh/ssh-nft.service
[Unit]
Description=ssh + nftables
After=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/ssh-nft.sh

